# pmi/agent/supervisord.conf
[supervisord]
nodaemon=true

[program:rabbitmq]
command=/opt/start-rabbitmq.sh

[program:agent]
command=/opt/start-agent.sh
